# git/etc/profile.d/aliases.sh
version https://git-lfs.github.com/spec/v1
oid sha256:8056f48ae61f2cba668c07b31c56229914912444ca6656205aaffce59b8e6b11
size 611
