# git/etc/profile.d/env.sh
version https://git-lfs.github.com/spec/v1
oid sha256:32c9a549ecb1c9e06622221dae98a000ac2ce9c2be97d5d0c4f80610af3dc55c
size 901
