# git/etc/profile.d/git-prompt.sh
version https://git-lfs.github.com/spec/v1
oid sha256:e2d0e4c58dca8ad7dba59df5db37a6ffb917f984edf67ae51d3f4c2db06e459b
size 1769
