# git/etc/profile.d/lang.sh
version https://git-lfs.github.com/spec/v1
oid sha256:6c0ca979c7b146b3750103b1296a399764f4e1b222ee091d3aa072b6da16c1a5
size 828
